# netsentinel-agent/deploy/macos/com.sounmu.netsentinel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sounmu.netsentinel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/netsentinel-agent</string>
    </array>

    <!-- dotenvy::dotenv() looks in CWD for .env, so point launchd at the
         protected config dir installed by deploy.sh. -->
    <key>WorkingDirectory</key>
    <string>/usr/local/etc/netsentinel</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/netsentinel/app.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/netsentinel/error.log</string>
</dict>
</plist>
